# config/gkd-unlock.service
[Unit]
Description=GNOME Keyring Unlock Service
After=graphical.target
Requires=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/gnome-keyring-unlock
Restart=on-failure

[Install]
WantedBy=default.target
